# base/etc/aldos/tmp.mount
[Mount]
Where=/tmp
What=tmpfs
Type=tmpfs
Options=rw,size=256M

[Install]
RequiredBy=local-fs.target
