# user/etc/aldos/user-files.service
[Service]
ExecStart=/etc/aldos/user-files.sh
