# user/etc/aldos/user-files.path
[Path]
PathChanged=/dev/disk/by-partlabel
