# user/etc/profile.d/start-desktop.sh
#!/bin/sh
if [[ $USER == "human" ]]
then
	startx
fi
